# templates/find.service.j2
[Unit]
Description=Find server.

[Service]
Type=forking
User={{ find_user }}
ExecStart=/usr/bin/findserver

[Install]
WantedBy=multi-user.target
